# ansible/roles/webmanager/files/metanyx-manager.service
[Unit]
Description=metanyx manager

[Service]
ExecStart=/opt/metanyx/manager/manage.py

[Install]
WantedBy=multi-user.target
