# ansible/roles/webmanager/files/torbone-manager.service
[Unit]
Description=TORbone manager

[Service]
ExecStart=/opt/torbone/manager/manage.py

[Install]
WantedBy=multi-user.target
